# scripts/player.service
[Unit]
Description=Player MP3 module
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=samos
WorkingDirectory=/home/samos/apps/Player
Environment=PLAYER_HOME=/home/samos/apps/Player
Environment=PLAYER_DISPLAY=ili9341
Environment=PLAYER_TFT_DC=25
Environment=PLAYER_TFT_RST=24
Environment=PLAYER_TFT_SPEED=24000000
Environment=PLAYER_TFT_REINIT_SECONDS=0
Environment=PLAYER_TFT_REINIT_ON_INPUT=0
Environment=PLAYER_AUDIO_DRIVER=alsa
Environment=PLAYER_AUDIO_DEVICE=alsa/default:CARD=Audio
Environment=PLAYER_ALSA_CARD=Audio
Environment=PLAYER_ALSA_CONTROL=Headphone
Environment=PLAYER_ALSA_VOLUME=85%
ExecStart=/usr/bin/python3 -m player.server --host 0.0.0.0 --port 8010
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
